# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# Gunicorn WSGI Server
[program:ai_trading_gunicorn]
command=/path/to/your/venv/bin/gunicorn --config gunicorn.conf.py ai_trading_engine.wsgi:application
directory=/path/to/your/project
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gunicorn.log
stderr_logfile=/var/log/supervisor/gunicorn_error.log
environment=DJANGO_SETTINGS_MODULE="ai_trading_engine.settings_production"
priority=100

# Celery Worker
[program:ai_trading_celery_worker]
command=/path/to/your/venv/bin/celery -A ai_trading_engine worker --loglevel=info --concurrency=4
directory=/path/to/your/project
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery_worker.log
stderr_logfile=/var/log/supervisor/celery_worker_error.log
environment=DJANGO_SETTINGS_MODULE="ai_trading_engine.settings_production"
priority=200

# Celery Beat (Scheduler)
[program:ai_trading_celery_beat]
command=/path/to/your/venv/bin/celery -A ai_trading_engine beat --loglevel=info
directory=/path/to/your/project
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery_beat.log
stderr_logfile=/var/log/supervisor/celery_beat_error.log
environment=DJANGO_SETTINGS_MODULE="ai_trading_engine.settings_production"
priority=300

# Celery Flower (Monitoring)
[program:ai_trading_celery_flower]
command=/path/to/your/venv/bin/celery -A ai_trading_engine flower --port=5555
directory=/path/to/your/project
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery_flower.log
stderr_logfile=/var/log/supervisor/celery_flower_error.log
environment=DJANGO_SETTINGS_MODULE="ai_trading_engine.settings_production"
priority=400

# Redis Server
[program:redis_server]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_error.log
priority=50

# Nginx Server
[program:nginx_server]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_error.log
priority=25

# Health Check Service
[program:health_check_service]
command=/path/to/your/venv/bin/python manage.py health_check
directory=/path/to/your/project
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/health_check.log
stderr_logfile=/var/log/supervisor/health_check_error.log
environment=DJANGO_SETTINGS_MODULE="ai_trading_engine.settings_production"
priority=500

# Performance Monitor
[program:performance_monitor]
command=/path/to/your/venv/bin/python manage.py performance_monitor
directory=/path/to/your/project
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/performance_monitor.log
stderr_logfile=/var/log/supervisor/performance_monitor_error.log
environment=DJANGO_SETTINGS_MODULE="ai_trading_engine.settings_production"
priority=600

# Backup Service
[program:backup_service]
command=/path/to/your/venv/bin/python manage.py backup_database
directory=/path/to/your/project
user=www-data
autostart=false
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backup.log
stderr_logfile=/var/log/supervisor/backup_error.log
environment=DJANGO_SETTINGS_MODULE="ai_trading_engine.settings_production"
priority=700

# Group configuration
[group:ai_trading_engine]
programs=ai_trading_gunicorn,ai_trading_celery_worker,ai_trading_celery_beat,ai_trading_celery_flower
priority=999

[group:infrastructure]
programs=redis_server,nginx_server
priority=999
